# src/main/resources/bin/upstart/etc/init/vidsync-server.conf
# vidsync server service

# Watches video directories and distributes to vidsync clients on your LAN.

description "Vidsync Server"
author "Geoff Bourne"

respawn

start on (local-filesystems and net-device-up IFACE!=lo)

env USER=xbmc
env HOME=/home/xbmc

setuid xbmc

script
   [ -f /etc/default/vidsync ] && . /etc/default/vidsync

   MY_IP=$(hostname -I|head -1)

   exec $JAVA_HOME/bin/java -Djava.rmi.server.hostname=$MY_IP -cp $VIDSYNC_JAR me.itzgeoff.vidsync.server.ServerMain
end script
